# lemonade-idle/systemd/lemonade-sampler.service
[Unit]
Description=Lemonade GPU utilization sampler
After=network.target

[Service]
ExecStart=/usr/bin/python3 %h/.local/share/lemonade-idle/sampler.py
Restart=always
RestartSec=10
Environment=LEMONADE_URL=http://192.168.0.125:3000
Environment=SAMPLE_INTERVAL_S=30
Environment=WINDOW_S=900

[Install]
WantedBy=default.target
